# conf.d/firefox.conf
[program:fluxbox]
command=firefox
autorestart=true
